# tools/frr.service.in
[Unit]
Description=FRRouting
Documentation=https://frrouting.readthedocs.io/en/latest/setup.html
After=networking.service systemd-sysctl.service csmgrd.service
OnFailure=heartbeat-failed@%n

[Service]
Nice=-5
Type=forking
NotifyAccess=all
StartLimitInterval=3m
StartLimitBurst=3
TimeoutSec=@TIMEOUT_MIN@m
WatchdogSec=60s
RestartSec=5
Restart=on-abnormal
LimitNOFILE=1024
LimitCORE=536870912
ExecStart=@CFG_SBIN@/frrinit.sh start
ExecStop=@CFG_SBIN@/frrinit.sh stop
ExecReload=@CFG_SBIN@/frrinit.sh reload


##security hardening
PrivateTmp=yes
RestrictNamespaces=~CLONE_NEWUTS
NotifyAccess=
PrivateDevices=true
ProtectHostname=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectProc=invisible

NoNewPrivileges=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
ProtectKernelLogs=true
ProtectKernelTunables=true

ProtectSystem=full
ProtectHome=yes
ProtectClock=yes
RestrictSUIDSGID=true
UMask=0077
CapabilityBoundingSet=~CAP_SYS_BOOT
CapabilityBoundingSet=~CAP_SYS_CHROOT


[Install]
WantedBy=multi-user.target
