# tools/frr.service.in
[Unit]
Description=FRRouting
Documentation=https://frrouting.readthedocs.io/en/latest/setup.html
Wants=network.target
After=network-pre.target systemd-sysctl.service csmgrd.service
Before=network.target
OnFailure=heartbeat-failed@%n

[Service]
Nice=-5
Type=forking
NotifyAccess=all
StartLimitInterval=3m
StartLimitBurst=3
TimeoutSec=@TIMEOUT_MIN@m
WatchdogSec=60s
RestartSec=5
Restart=on-abnormal
LimitNOFILE=1024
LimitCORE=536870912
ExecStart=@CFG_SBIN@/frrinit.sh start
ExecStop=@CFG_SBIN@/frrinit.sh stop
ExecReload=@CFG_SBIN@/frrinit.sh reload
Environment=ASAN_OPTIONS=print_summary=1:detect_deadlocks=1:check_printf=1:detect_leaks=1:print_full_thread_history=1:halt_on_error=0:log_path=/tmp/asan.log:log_exe_name=1
LimitAS=infinity
LimitDATA=infinity

[Install]
WantedBy=multi-user.target
